# event-gateway/supervisord.recalc.conf
[program:recalc_worker]
command=php /app/worker/recalc_worker.php
autostart=true
autorestart=true
stderr_logfile=/var/log/recalc.err.log
stdout_logfile=/var/log/recalc.out.log
