# event-gateway/supervisord.monitor.conf
[program:monitor_trigger]
command=php /app/monitor/monitor_trigger.php
autostart=true
autorestart=true
stderr_logfile=/var/log/monitor.err.log
stdout_logfile=/var/log/monitor.out.log
